# build/debian/systemd.jaspy/jaspy-nexus.service
[Unit]
Description=jaspy-nexus
After=network.target postgresql.service

[Service]
Type=simple
User=jaspy
Group=jaspy
Environment=JASPY_DB_URL=postgresql:///jaspy
Environment=JASPY_MQTT_SERVER=127.0.0.1
Environment=ROCKET_DATABASES={jaspy_db={url="postgres:///jaspy"}}
ExecStart=/usr/lib/jaspy/jaspy-nexus

[Install]
WantedBy=multi-user.target

